# supervisord.conf
[supervisord]
nodaemon=true

[program:cellect]
user=root
command=/cellect_panoptes/start
directory=/cellect_panoptes
autorestart=true
